# voice-daemon.service
[Unit]
Description=Loki voice inbox daemon — answers Iðunn's voice messages
After=network.target

[Service]
User=bmo
WorkingDirectory=/home/bmo/pucky
ExecStart=/usr/bin/python3 /home/bmo/pucky/voice_daemon.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
